# deploy/supervisord.conf
; deploy/supervisord.conf
; ────────────────────────
; supervisord manages two processes inside the HF Spaces container.
;
; WHY supervisord AND NOT TWO SEPARATE CONTAINERS?
;
;   HF Spaces (free tier) only gives you one container and one exposed port.
;   supervisord is the standard way to run multiple processes inside a single
;   container when you can't use docker-compose or kubernetes.
;
; PROCESS ORDER:
;
;   1. uvicorn starts first (priority=10). The FastAPI server loads the
;      DistilBERT model into memory — this takes ~1-2 seconds.
;
;   2. streamlit starts after a 3-second delay (startsecs=3) to give uvicorn
;      time to finish model loading before the UI tries to call /health.
;
; LOGGING:
;
;   Both processes log to /tmp/ (writable by our non-root user).
;   HF Spaces streams these logs in the Space's "Logs" tab.

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid
loglevel=info

; ── FastAPI server (internal, port 8000) ─────────────────────────────────────
[program:api]
command=python -m uvicorn api.app:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
priority=10
startsecs=5
stdout_logfile=/tmp/api_stdout.log
stderr_logfile=/tmp/api_stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

; ── Streamlit UI (public, port 7860) ─────────────────────────────────────────
[program:streamlit]
command=python -m streamlit run app/streamlit_app.py
    --server.port 7860
    --server.address 0.0.0.0
    --server.headless true
    --server.enableCORS false
    --server.enableXsrfProtection false
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=3
; Point Streamlit at the internal API server
environment=API_URL="http://localhost:8000"
stdout_logfile=/tmp/streamlit_stdout.log
stderr_logfile=/tmp/streamlit_stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
